# approot/supervisord.conf
[program:node]
command = node tripsaver.js
directory = /home/dotcloud/current